# build/config/systemd/metadatastorer.service
# Systemd service at /etc/systemd/system/rsloadmanager.service

[Unit]
After=cloud-final.service
Wants=cloud-final.service
StartLimitInterval=0

[Service]
ExecStart=/opt/science/rs_ingester/bin/run_metadatastorer.sh
Restart=always
RestartSec=5
TimeoutStopSec=10

[Install]
WantedBy=cloud-init.target
