# build/config/upstart/dbstorer.conf
# Upstart script at /etc/init/redshiftIngester.conf

start on runlevel [2345]
stop on runlevel [06]

kill signal SIGINT
kill timeout 600

respawn
respawn limit 10 5

script
  exec /opt/science/redshiftIngester/bin/run_dbstorer.sh
  emit storer_running
end script
