# build/config/upstart/metadatastorer.conf
# Upstart script at /etc/init/metadatastorer.conf

start on runlevel [2345] and stopped cloud-final
stop on runlevel [06]

kill signal SIGINT
kill timeout 600

respawn
respawn limit 10 5

script
  exec /opt/science/rs_ingester/bin/run_metadatastorer.sh
  emit storer_running
end script
